# oaic/srsenb.service
[Unit]
Description=SRS Enhanced NodeB

[Service]
ExecStart=/usr/local/bin/start-srsenb
Type=exec

[Install]
WantedBy=default.target
