# oaic/srsepc.service
[Unit]
Description=SRS Enhanced Packet Core

[Service]
ExecStart=/usr/local/bin/srsepc
Type=exec

[Install]
WantedBy=default.target
